# root/etc/s6-overlay/s6-rc.d/myapp/run
#!/usr/bin/with-contenv bash

echo "Hello World" > /var/log/myapp.log